# scripts/lua-server.service
[Unit]
Description=Lua scraper    

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/usr/bin/scraper.sh

[Install]
WantedBy=default.target
